# dotfiles/.shell.d/config.d/cron/cron.daily/clamscan
#!/bin/sh
clamscan -r --quiet --bell -i /home/$USER/Downloads
